# wsgi/uwsgi.conf
description "uWSGI Emperor"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /home/poyzer/envs/ordergroove/bin/uwsgi --max-fd 65000 --emperor /etc/uwsgi/vassals/
